# launchd/com.pugs.syncagent.updater.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.pugs.syncagent.updater</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>__AGENT_ROOT__/update.sh</string>
  </array>

  <!-- Check for updates every 10 minutes -->
  <key>StartInterval</key>
  <integer>600</integer>
  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>__AGENT_ROOT__/updater.log</string>
  <key>StandardErrorPath</key>
  <string>__AGENT_ROOT__/updater.error.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
</dict>
</plist>
